# openclaw-compose.service
[Unit]
Description=OpenClaw + ClawPanel Docker Compose Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
ExecStart=__DOCKER_BIN__ compose up -d --no-build --remove-orphans
RemainAfterExit=yes
Restart=on-failure
RestartSec=10s
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
